# efinder-ensure-ap.service
[Unit]
Description=Ensure eFinder Wi-Fi AP is active
After=NetworkManager.service network.target
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/efinder-ensure-ap.sh
RemainAfterExit=yes
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
