# image-config/efinder-usb-gadget.service
[Unit]
Description=eFinder USB serial ACM gadget (configfs)
After=systemd-modules-load.service sys.mount
Wants=systemd-modules-load.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/efinder-gadget-connect
ExecStartPost=-/bin/systemctl start serial-getty@ttyGS0.service
Restart=on-failure
RestartSec=15s
StartLimitIntervalSec=120s
StartLimitBurst=4

[Install]
WantedBy=multi-user.target
