# image-config/efinder-firstboot.service
[Unit]
Description=eFinder boot-time configuration
# Need NetworkManager up so we can nmcli con add the AP profile and check/fix
# the AP autoconnect settings.  Don't wait for network-online -- on first
# boot we *create* the network, so network-online.target may never fire.
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/efinder-firstboot.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
# First boot can take 30-60 seconds (waiting for NM, hardware checks).
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
